# ubuntu/etc/monit.conf
description "A utility for monitoring and managing daemons"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
exec /usr/sbin/monit -Ic /etc/monit/monitrc 
